# files/network-assistant-client.timer
[Unit]
Description=Scheduled Run Network Assistant Client

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
